# examples/systemd/mage-s3.service.example
[Unit]
Description=Mage image processing service (S3 storage)
Documentation=https://github.com/sashko-guz/mage
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mage
Group=mage
WorkingDirectory=/opt/mage

# Storage
Environment=STORAGE_DRIVER=s3
Environment=S3_BUCKET=your-bucket-name
Environment=S3_REGION=us-east-1
# S3_ACCESS_KEY and S3_SECRET_KEY should be in .env file
# For S3-compatible storage (MinIO):
# Environment=S3_BASE_URL=https://minio.example.com
# Environment=S3_USE_PATH_STYLE=true

# Server
Environment=PORT=8080
Environment=LOG_LEVEL=info

# S3 HTTP tuning (optional)
Environment=S3_MAX_IDLE_CONNS=150
Environment=S3_MAX_IDLE_CONNS_PER_HOST=150

# Cache (recommended for S3 to reduce latency)
Environment=SOURCE_MEMORY_CACHE_ENABLED=true
Environment=SOURCE_MEMORY_CACHE_MAX_SIZE_MB=512
Environment=SOURCE_DISK_CACHE_ENABLED=true
Environment=SOURCE_DISK_CACHE_DIR=/var/cache/mage/sources
Environment=SOURCE_DISK_CACHE_MAX_SIZE_MB=10240
Environment=THUMB_MEMORY_CACHE_ENABLED=true
Environment=THUMB_MEMORY_CACHE_MAX_SIZE_MB=1024
Environment=THUMB_DISK_CACHE_ENABLED=true
Environment=THUMB_DISK_CACHE_DIR=/var/cache/mage/thumbs
Environment=THUMB_DISK_CACHE_MAX_SIZE_MB=20480

# Load additional config from file (secrets, overrides)
EnvironmentFile=-/opt/mage/.env

ExecStart=/usr/local/bin/mage
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/cache/mage

# Resource limits
LimitNOFILE=65535
MemoryMax=4G

[Install]
WantedBy=multi-user.target
