# examples/systemd/mage-local.service.example
[Unit]
Description=Mage image processing service (local storage)
Documentation=https://github.com/sashko-guz/mage
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mage
Group=mage
WorkingDirectory=/opt/mage

# Storage
Environment=STORAGE_DRIVER=local
Environment=STORAGE_ROOT=/var/lib/mage/images

# Server
Environment=PORT=8080
Environment=LOG_LEVEL=info

# Cache (optional - enable as needed)
Environment=SOURCE_MEMORY_CACHE_ENABLED=false
Environment=SOURCE_DISK_CACHE_ENABLED=false
Environment=SOURCE_DISK_CACHE_DIR=/var/cache/mage/sources
Environment=THUMB_MEMORY_CACHE_ENABLED=true
Environment=THUMB_MEMORY_CACHE_MAX_SIZE_MB=512
Environment=THUMB_DISK_CACHE_ENABLED=true
Environment=THUMB_DISK_CACHE_DIR=/var/cache/mage/thumbs
Environment=THUMB_DISK_CACHE_MAX_SIZE_MB=10240

# Load additional config from file (secrets, overrides)
EnvironmentFile=-/opt/mage/.env

ExecStart=/usr/local/bin/mage
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/cache/mage

# Resource limits
LimitNOFILE=65535
MemoryMax=2G

[Install]
WantedBy=multi-user.target
